# files/etc/systemd/system/plexmediaserver.service.d/override.conf
[Service]
Environment="PLEX_MEDIA_SERVER_USE_SYSLOG=true"
